# scripts/linux/jarvis-ui.service
[Unit]
Description=Jarvis UI (Chromium app-mode)
After=jarvis-backend.service graphical-session.target
Requires=jarvis-backend.service
# Tie lifetime to the Wayland session: stops on logout, restarts on login.
PartOf=graphical-session.target

[Service]
Type=simple
ExecStartPre=/bin/sleep 2
ExecStart=/usr/bin/chromium \
  --app=http://localhost:8788 \
  --no-default-browser-check \
  --disable-infobars \
  --autoplay-policy=no-user-gesture-required \
  --use-fake-ui-for-media-stream=false \
  --ozone-platform-hint=auto \
  --enable-features=UseOzonePlatform
Restart=on-failure
RestartSec=3
# DISPLAY/WAYLAND_DISPLAY come from the session environment — UWSM imports it
# into the systemd user manager; --ozone-platform-hint=auto picks the backend.

[Install]
WantedBy=graphical-session.target
