# scripts/linux/jarvis-tts.service
[Unit]
Description=Jarvis TTS Service (edge-tts / XTTS)
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/jarvis-linux/backend/voice/python
ExecStart=%h/jarvis-linux/backend/voice/python/.venv/bin/python edge_tts_service.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
