# scripts/linux/jarvis-wake.service
[Unit]
Description=Jarvis Wake Word Daemon (openWakeWord)
After=jarvis-backend.service

[Service]
Type=simple
WorkingDirectory=%h/jarvis-linux/backend/voice/python
ExecStart=%h/jarvis-linux/backend/voice/python/.venv/bin/python wake_service.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
