# scripts/linux/jarvis-jarvisbot.service
[Unit]
Description=Jarvis Telegram Bot (JarvisBot)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/jarvis-linux
ExecStart=/usr/bin/node backend/src/services/jarvisBot.js
Restart=on-failure
RestartSec=10
# Exit 0 (no token configured) does NOT restart — intentional.
RestartPreventExitStatus=0

[Install]
WantedBy=default.target
